# debian/xscreensaver.service
[Unit]
Description=XScreenSaver

[Service]
ExecStart=xscreensaver

[Install]
WantedBy=default.target
